# deployment/systemd/recipes-frontend.service
[Unit]
Description=Recipes Application Frontend (Node/Vite Preview)
After=network.target

[Service]
Type=simple
# Change these to match your setup
User=recipes
Group=recipes
WorkingDirectory=/opt/recipes/frontend

# Serve the built frontend using vite preview
# For production, consider using nginx instead
ExecStart=/usr/bin/npm run preview -- --host 0.0.0.0 --port 5173

# Restart policy
Restart=on-failure
RestartSec=5s

# Security settings
NoNewPrivileges=true
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=recipes-frontend

[Install]
WantedBy=multi-user.target
